# etc/systemd/cf-ddns.timer
[Unit]
Description=cf-ddns timer

[Timer]
OnUnitActiveSec=5
Unit=cf-ddns.service

[Install]
WantedBy=multi-user.target