# etc/systemd/cf-ddns.service
[Unit]
Description=cf ddns service

[Service]
ExecStart=cf_ddns.py